# watchping/files/etc/uci-defaults/watching
#!/bin/sh
#Copyright ©  dungtd8x@gmail.com

#/etc/init.d/watchping disable
#/etc/init.d/watchping stop
