# advance/files/initd_network_reload
#!/bin/sh /etc/rc.common

START=99
STOP=90

boot() {
    /etc/init.d/network restart
	/sbin/wifimedia/controller.sh _boot
}
