# advance/files/initd_network_reload
#!/bin/sh /etc/rc.common

START=90
STOP=12

boot() {
	/etc/init.d/network restart
	if [ "$(uci -q get wifimedia.@lte[0].4glte)" -eq 0 ] || [ "$(uci -q get wifimedia.@lte[0].4glte)" = "" ];then
		echo 0 >/sys/class/gpio/power_usb3/value
	fi
	/sbin/wifimedia/controller.sh _boot
}

