# activation/files/initd_wifimedia_check
#!/bin/sh /etc/rc.common

START=30
STOP=105

boot() {
	/sbin/wifimedia/controller.sh _lic >/dev/null 2>&1
}
